# deploy/systemd/mockzure.service
[Unit]
Description=Mockzure Service (Mock Azure for testing)
Documentation=https://github.com/yourcloudtools/sandman
After=network.target

[Service]
Type=simple
User=sandman
Group=sandman
WorkingDirectory=/var/lib/mockzure
ExecStart=/usr/bin/mockzure
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mockzure

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PORT=8090"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target

